# deploy/docker/templates/cron.d/renew-certificate
* * * * 0 root exec /opt/bizBrainz/renew-certificate.sh
